# sv/parity/run
#!/bin/bash

source /etc/envvars

DATA=${DATA="/data"}
mkdir -p $DATA

exec 2>&1
eval "exec /parity --db-path=$DATA $ARGS"
